# configuration/roles/lego/templates/lego.service.j2
[Unit]
Description=Issue or renew ACME certificates via lego
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile={{ lego_env_file }}
ExecStart=/usr/local/bin/lego-run.sh
